# systemd/bluealsa.service.d/override.conf
# =============================================================================
# bt2iap — systemd drop-in override for bluealsa.service
# -----------------------------------------------------------------------------
# Path on Pi: /etc/systemd/system/bluealsa.service.d/override.conf
#
# Purpose:
#   The bluealsa package ships a unit that launches `/usr/bin/bluealsa` with
#   default profiles. We want the daemon to advertise ONLY A2DP sink (we are
#   the receiver; the phone is the source), not HFP / HSP, not A2DP source.
#
#   systemd drop-ins merge with the shipped unit. Lines set here override the
#   shipped unit's matching lines.
#
# Apply:
#   sudo install -d /etc/systemd/system/bluealsa.service.d
#   sudo install -m 0644 systemd/bluealsa.service.d/override.conf \
#       /etc/systemd/system/bluealsa.service.d/override.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart bluealsa
#
# Verify:
#   systemctl show bluealsa -p ExecStart --value | head -n 1
#   # → should list our /usr/bin/bluealsa -p a2dp-sink ... argv
#   journalctl -u bluealsa -n 50
#   # → should show "BlueALSA starting" and no HFP/HSP init lines
# =============================================================================

[Service]

# Clear the shipped ExecStart. This empty assignment is the canonical systemd
# idiom for "discard inherited ExecStart before setting a new one". Without
# this line, systemd would append (and reject) a second ExecStart directive.
ExecStart=

# Restrict BlueALSA to A2DP sink only. Flags:
#   -p a2dp-sink    enable only the A2DP sink profile (phone → Pi audio)
#   --keep-alive=5  keep the PCM alive for 5s after the stream stops, so a
#                   short pause on the phone does not tear down the ALSA
#                   pipeline (the loopback + g_ipod_audio chain downstream
#                   is cheap to keep warm and expensive to rebuild).
ExecStart=/usr/bin/bluealsa -p a2dp-sink --keep-alive=5

# Raise LIBASOUND_DEBUG to 1 to force verbose ALSA logging when triaging
# broken audio; default 0 keeps journalctl quiet in normal operation.
Environment=LIBASOUND_DEBUG=0
